# deploy/systemd/braiins-monitor.service
[Unit]
Description=Braiins bot and service monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
EnvironmentFile=/etc/braiins-monitor/env
ExecStart=/usr/local/bin/braiins-monitor run
WorkingDirectory=/opt/braiins-monitor
Nice=10
